# N1/Lede/files/etc/uci-defaults/99-docker
#!/bin/sh
/etc/init.d/dockerd enable
exit 0
